# systemd/division-overtime-employee-consistency.timer
[Unit]
Description=Record Division Overtime employee consistency history daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=division-overtime-employee-consistency.service

[Install]
WantedBy=timers.target
