# systemd/division-overtime-threshold.service
[Unit]
Description=Division Overtime threshold notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/division-overtime
EnvironmentFile=/home/pi/division-overtime/.env
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/division-overtime/.venv/bin/division-overtime --root /home/pi/division-overtime run threshold
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/division-overtime/var
