# systemd/division-overtime-threshold.timer
[Unit]
Description=Run Division Overtime threshold check on weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 10:30:00
Persistent=true
RandomizedDelaySec=30
Unit=division-overtime-threshold.service

[Install]
WantedBy=timers.target
