# systemd/division-overtime-web.service
[Unit]
Description=division-overtime Web administration
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/division-overtime
EnvironmentFile=/home/pi/division-overtime/.env
ExecStart=/home/pi/division-overtime/.venv/bin/division-overtime-web
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/pi/division-overtime/var /home/pi/division-overtime/data

[Install]
WantedBy=multi-user.target
